# configs/openrc/netinfo_nconnect
#!/sbin/openrc-run
# https://github.com/OpenRC/openrc/blob/master/service-script-guide.md
# env
export N_ID="xxx"
export N_ENDPOINT="xxx"
export N_INTERVAL="5m"

command="/usr/local/bin/netinfo"
command_args="send nconnect -id=${N_ID} -endpoint=${N_ENDPOINT} -allow_insecure -interval=${N_INTERVAL}"
command_background=true
command_user="root"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
  need net
  use dns logger netmount
  after firewall
}
